# .config/profile.d/02-programs.sh
export ALTERNATE_EDITOR=""
export EDITOR='emacsclient -c'
export VISUAL="${EDITOR} -n"
export TERMINAL="foot"
export BROWSER="firefox"
export GRIM_DEFAULT_DIR="$HOME/pics/screen/"
export BEMENU_OPTS="-i -l 10"
