# .config/profile.d/03-guix.sh
export GUIX_LOCPATH=$HOME/.guix-profile/lib/locale
