# .config/profile.d/path.sh
LOCAL=$HOME/.local
SHARE=$LOCAL/share

export PATH="$PATH:$LOCAL/bin"
