# .config/profile.d/path.sh
LOCAL=$HOME/.local
SHARE=$LOCAL/share

export PATH="$LOCAL/bin:$PATH"
