# .config/profile.d/programs.sh
export EDITOR="emacsclient -a \"\" -c"
export VISUAL="$EDITOR -n"
export TERMINAL="$VISUAL -e '(eshell t)'"
export QTERM="$VISUAL -e '(eshell)'"
export BROWSER="brave"
